# multipathd/96-multipath-late.rules
# event to be catched by multipathd
RUN+="socket:/org/kernel/dm/multipath_event"

